# services/pyload/run
#!/bin/bash

echo "Starting Pyload"
exec /sbin/setuser abc /opt/pyload/pyLoadCore.py
